# plays/roles/workstation/files/autostart/shutter.desktop
[Desktop Entry]
Name=Shutter
Exec=shutter
Icon=shutter
Terminal=false
Type=Application
