# plays/roles/workstation/files/autostart/shutter.desktop
[Desktop Entry]
Name=Shutter
Exec=shutter --min_at_startup
Icon=shutter
Terminal=false
Type=Application
